# com.aronhuang.voice-capture.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.aronhuang.voice-capture</string>
    <key>ProgramArguments</key>
    <array>
        <string>/Users/aronhuang/.venvs/whisper/bin/python3</string>
        <string>/Users/aronhuang/Developer/Kai/Voice-Capture/watcher.py</string>
    </array>
    <key>WatchPaths</key>
    <array>
        <string>/Users/aronhuang/Library/Group Containers/group.com.apple.VoiceMemos.shared/Recordings</string>
    </array>
    <key>EnvironmentVariables</key>
    <dict>
        <key>PATH</key>
        <string>/opt/homebrew/bin:/usr/local/bin:/usr/bin:/bin</string>
    </dict>
    <key>StandardOutPath</key>
    <string>/tmp/voice-capture.stdout.log</string>
    <key>StandardErrorPath</key>
    <string>/tmp/voice-capture.stderr.log</string>
</dict>
</plist>
